# config_files/supervisord.conf
[supervisord]
nodaemon=true
user=root ; as we are inside the Docker container

[program:redis-server]
priority=1
command=redis-server
redirect_stderr=true
stdout_logfile=/process_logs/redis-server.log

[program:worker]
priority=2
directory=/repo
command=rq worker
process_name=%(program_name)s-%(process_num)s
numprocs=15
redirect_stderr=true
stdout_logfile=/process_logs/%(program_name)s-%(process_num)s.log

[program:uwsgi]
priority=3
command=uwsgi --ini /repo/config_files/cagecat.ini --safe-pidfile /tmp/uwsgi-master.pid
redirect_stderr=true
stdout_logfile=/process_logs/uwsgi.log

[program:nginx]
priority=4
command=nginx
redirect_stderr=true
stdout_logfile=/process_logs/nginx.log

[program:cron]
priority=5
command=cron
; no logging as it is daemon
